# unlinksnoop.service
[Service]
ExecStart=/usr/bin/unlinksnoop
KillMode=none

[Install]
WantedBy=basic.target